# deploy/typeweaver.service
[Unit]
Description=TypeWeaver Typography Robustness Lab
After=network.target

[Service]
Type=simple
User=typeweaver
Group=typeweaver
WorkingDirectory=/opt/typeweaver
ExecStart=/opt/typeweaver/typeweaver serve --registry-root /opt/typeweaver/.typeweaver --host 0.0.0.0 --port 3500
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
